# deploy/ottapi-v2.service
[Unit]
Description=OTTAPI v2 - IPTV catalog backend
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/ottapi
Environment="PATH=/home/ottapi/.venv/bin:/usr/bin:/usr/local/bin"
ExecStart=/home/ottapi/.venv/bin/uvicorn src.main:app --host 127.0.0.1 --port 8011 --workers 1
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ottapi-v2

[Install]
WantedBy=multi-user.target
